# ansible/roles/fileserver/files/health-check/cron.d/health-check
MAILTO=''
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 * * * * root /root/health-check/health-check.sh {{ ansible_hostname }} {{ smtp_email }} {{ smtp_notification_email }}
